# IotController/deploy/smartsecurity.service
[Unit]
Description=Smart Security System - FSM Controller
Documentation=https://github.com/NotJar404/SmartSecuritySystem
After=network.target webapp.service
Wants=network-online.target webapp.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/SmartSecuritySystem/IotController
ExecStart=/home/pi/SmartSecuritySystem/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# GPIO permissions
SupplementaryGroups=gpio spi i2c video

[Install]
WantedBy=multi-user.target
